# supervisor.ini
[program:markovifybot]
command=/usr/bin/python3 /home/wodim/markovifybot/telegram_bot.py
directory=/home/wodim/markovifybot
user=wodim
autorestart=true
autostart=true
